# package/opengw/root/etc/init.d/opengw
#!/bin/sh /etc/rc.common
START=1999
start(){
        /usr/share/opengw/opengw >/tmp/opengw.log &
}

stop(){
       killall opengw &
}

restart(){
        killall opengw && /usr/share/opengw &
}
